# ipakeyapp.conf
[program:ipakeyapp]
command=python /var/www/ipakeyapp/app.py
redirect_stderr=true
stdout_logfile=/var/www/ipakeyapp/app.log